# docs/ugv-beast/setup/ugv_tools_api/deploy/ugv-ears.service
[Unit]
Description=UGV Ears (mic -> wake-word -> STT -> /chat)
# Full chain: ears needs tools_api (for ROS/camera), voice (for auto-speak),
# AND ugv-er (mission endpoint receives transcripts).
Requires=ugv-tools-api.service ugv-voice.service ugv-er.service
After=ugv-tools-api.service ugv-voice.service ugv-er.service
# Keep retrying — USB mic enumeration can lag ~20s after boot.
StartLimitBurst=5
StartLimitIntervalSec=60

[Service]
Type=simple
Restart=on-failure
RestartSec=5
EnvironmentFile=/home/jetson/ugv_ws_waveshare/ugv_tools_api/deploy/ears.env
EnvironmentFile=/home/jetson/ugv_ws_waveshare/ugv_tools_api/deploy/er.env
# /bin/bash so a missing +x does not fail the service.
ExecStart=/usr/bin/docker exec -e BLACKBOX_OPERATOR -e MIC_DEVICE_HINT -e WAKEWORD_MODEL -e WAKEWORD_THRESHOLD -e VAD_AGGRESSIVENESS -e MUTE_FLAG_PATH -e OPENAI_API_KEY -e OPENAI_STT_MODEL -e ER_URL -e SPEAK_URL -e OPERATOR -e PULSE_SERVER -e PULSE_SOURCE -e PULSE_SINK -e UGV_ALSA_DEVICE -e EARS_CHIME -e EARS_CHIME_VOLUME -e WAKE_CHIME_PATH ugv_waveshare /bin/bash /home/ws/ugv_ws/ugv_tools_api/deploy/start_ears.sh
ExecStop=/usr/bin/docker exec ugv_waveshare pkill -f "ugv_tools_api.voice.ears"
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
